# scripts/ai-research-daily.service
[Unit]
Description=AI Research LinkedIn Generator - Daily
After=network.target

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/home/ubuntu/AI_EMPLOYEE_APP
Environment="PYTHONPATH=/home/ubuntu/AI_EMPLOYEE_APP"
Environment="GLM_API_KEY=${GLM_API_KEY}"
Environment="GLM_API_URL=https://api.z.ai/api/coding/paas/v4"
Environment="USE_CDP=false"
Environment="HEADLESS=true"
ExecStart=/usr/bin/python3 .claude/skills/research-linkedin-generator/scripts/research.py --daily
StandardOutput=append:/home/ubuntu/AI_EMPLOYEE_APP/logs/research-service.log
StandardError=append:/home/ubuntu/AI_EMPLOYEE_APP/logs/research-service.log

[Install]
WantedBy=multi-user.target
