# scripts/ai-research-daily.timer
[Unit]
Description=Daily AI Research LinkedIn Generator
Requires=ai-research-daily.service

[Timer]
OnCalendar=daily
OnCalendar=09:00
Persistent=true

[Install]
WantedBy=timers.target
